# swupd-probe.path
[Unit]
Description=swupd telemetry probe
ConditionPathExists=/etc/telemetrics/opt-in
Requires=telemprobd.socket
After=telemprobd.socket

[Path]
DirectoryNotEmpty=/var/lib/swupd/telemetry

[Install]
WantedBy=multi-user.target
