# swupd-probe.service
[Unit]
Description=swupd telemetry probe
ConditionPathExists=/etc/telemetrics/opt-in

[Service]
Type=simple
ExecStart=/usr/bin/swupd-probe
